# 0x1A-application_server/gunicorn.service
# /etc/systemd/system/gunicorn.service

[Unit]
Description=Gunicorn instance to serve AirBnB_clone_v4
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/path/to/AirBnB_clone_v4
ExecStart=/path/to/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:5003 --error-logfile /tmp/airbnb-error.log --access-logfile /tmp/airbnb-access.log web_dynamic.2-hbnb:app

[Install]
WantedBy=multi-user.target